# grafana/supervisord.conf
[program:http]
command=/opt/grafana/bin/grafana-server
numprocs=1
autostart=true
autorestart=true
